# systemd/dmujeres.service
[Unit]
Description=Dmujeres-Traccar Go Backend Service
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/Dmujeres-Traccar/backend
ExecStart=/home/Dmujeres-Traccar/backend/dmujeres-backend
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=dmujeres-traccar
Environment="DATABASE_URL=postgres://traccar:traccar_pass_556@127.0.0.1:5432/traccar?sslmode=disable"
Environment="PORT=8082"
Environment="REDIS_URL=localhost:6379"
Environment="TZ=America/Guayaquil"

[Install]
WantedBy=multi-user.target
